# docker/dragon-backup.timer
[Unit]
Description=Nightly Dragon Idea Engine Backup

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
Unit=dragon-backup.service

[Install]
WantedBy=timers.target
