# docker/dragon-alert-check.service
[Unit]
Description=Dragon Idea Engine Alert Check
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__REPO_DIR__
ExecStart=/bin/bash -lc 'cd "__REPO_DIR__" && output="$(./scripts/pi-alert-check.sh 2>&1)" && printf "%s\n" "$output" || { status=$?; printf "%s\n" "$output"; ALERT_MESSAGE="$output" ./scripts/pi-alert-notify.sh || true; exit $status; }'
