# docker/dragon-update.service
[Unit]
Description=Dragon Idea Engine Update Job
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__REPO_DIR__
ExecStart=/bin/bash -lc 'cd "__REPO_DIR__" && ./scripts/update-pi.sh'
